# templates/dev.postgresql.plist.erb
<?xml version="1.0" encoding="UTF-8"?>

<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.postgresql</string>

    <key>Boxen</key>
    <dict>
      <key>Description</key>
      <string>A tolerably okay relational database.</string>

      <key>Host</key>
      <string><%= @host %></string>

      <key>Port</key>
      <integer><%= @port %></integer>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string><%= @executable %></string>
      <string>-p</string>
      <string><%= @port %></string>
      <string>-D</string>
      <string><%= @datadir %></string>
      <string>-r</string>
      <string><%= @logdir %>/server.log</string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <true />

    <key>UserName</key>
    <string><%= scope.lookupvar "::boxen_user" %></string>

    <key>WorkingDirectory</key>
    <string><%= @datadir %></string>

    <key>StandardErrorPath</key>
    <string><%= @logdir %>/console.log</string>

    <key>StandardOutPath</key>
    <string><%= @logdir %>/console.log</string>
  </dict>
</plist>
